# ansible/roles/parasight_webserver/templates/systemd/gunicorn-parasight.socket.j2
[Unit]
Description=parasight gunicorn socket

[Socket]
ListenStream=/tmp/gunicorn-parasight.sock
SocketMode=666

[Install]
WantedBy=sockets.target
